# crontab
05 * * * * /media/sf_PROJET/addInformationsOptions.sh -e
